# debian/gnunet-server.prerm
#!/bin/sh

set -e

if [ -x /etc/init.d/gnunet-daemon ]
then
	if [ -x /usr/sbin/invoke-rc.d ]
	then
		invoke-rc.d gnunet-daemon stop || exit 0
	else
		/etc/init.d/gnunet-daemon stop || exit 0
	fi
fi

#DEBHELPER#

exit 0
