# debian/gnunet-daemon.postrm
#!/bin/sh -e
if [ "$1" = "purge" ] ; then
	update-rc.d -f gnunet-daemon remove > /dev/null 2>&1 || exit 0
fi

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0

