# debian/gnunet.prerm
#!/bin/sh
set -e

if [ -x "/etc/init.d/gnunet" ]; then
	if [ -x /usr/sbin/invoke-rc.d ] ; then
		invoke-rc.d gnunet stop || exit 0
	else
		/etc/init.d/gnunet stop || exit 0
	fi
fi

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0

