# debian/gnunet-server.prerm
#!/bin/sh

set -e

if [ -x /etc/init.d/gnunet-server ]
then
	if [ -x /usr/sbin/invoke-rc.d ]
	then
		invoke-rc.d gnunet-server stop || exit "${?}"
	else
		/etc/init.d/gnunet-server stop || exit "${?}"
	fi
fi

#DEBHELPER#

exit 0
