# debian/gnunet-server.postrm
#!/bin/sh

set -e

case "${1}" in
	purge)
		update-rc.d -f gnunet-server remove > /dev/null 2>&1 || true

		rm -rf /var/log/gnunetd
		rm -rf /var/run/gnunetd
		;;

	remove|upgrade|failed-upgrade|abort-install|abort-upgrade|disappear)

		;;

	*)
		echo "postrm called with unknown argument \`${1}'" >&2
		exit 1
		;;
esac

#DEBHELPER#

exit 0
